# addons/dynamic-mount/base/inotify-fluid-config.ini
[program:inotify-fluid-config]
command=/usr/local/bin/inotify.sh
redirect_stderr=true
stdout_logfile=/proc/1/fd/1
stdout_logfile_maxbytes=0
autorestart=true
startretries=9999